# conf/systemd/rc-local.timer
[Unit]
Description=run rc-local later

[Timer]
OnStartupSec=60
Unit=delay.service

[Install]
Alias=rc-local.timer
WantedBy=multi-user.target
