# conf/services/rc-local.timer
[Unit]
Description=run rc-local later

[Timer]
OnStartupSec=60
Unit=rc-local.service

[Install]
Alias=rc-local.timer
WantedBy=multi-user.target
